# config/rust-simple-tunnel.service
[Unit]
Description=simple rust traffic filter
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=$BIN_PATH -c $CONFIG_PATH --forward-traffic $ROUTE

[Install]
WantedBy=multi-user.target
